# mt5_supervisor.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

; ── 1. Virtual framebuffer — gives MT5 a display without a real screen ────────
; 800x600x16 is the minimum. 24-bit colour wastes memory; we don't render UI.
[program:xvfb]
command=Xvfb :99 -screen 0 800x600x16 -ac -nolisten tcp
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/null
stderr_logfile=/dev/null

; ── 2. MetaTrader 5 terminal ──────────────────────────────────────────────────
; /portable keeps all data under the install dir (inside the volume).
; WINEDEBUG=-all suppresses the flood of Wine debug messages.
[program:mt5]
command=/bin/bash -c "DISPLAY=:99 WINEDEBUG=-all wine 'C:\\Program Files\\MetaTrader 5\\terminal64.exe' /portable"
autostart=true
autorestart=true
priority=20
startretries=5
startsecs=10
stdout_logfile=/dev/null
stderr_logfile=/dev/null

; ── 3. mt5linux RPyC server ───────────────────────────────────────────────────
; Starts the Windows Python RPyC server inside Wine. The signal_platform
; connects to this on port 8812 from the Linux side via mt5linux.
; startsecs=25 waits for MT5 to fully initialise before binding port 8812.
[program:mt5_rpyc]
command=/bin/bash -c "DISPLAY=:99 WINEDEBUG=-all wine python.exe -m mt5linux --host 0.0.0.0 --port 8812"
autostart=true
autorestart=true
priority=30
startretries=10
startsecs=25
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
