# python/copy_trading/supervisor/copyengine.conf
; ─────────────────────────────────────────────────────────────────────────────
; TradeSync Copy Engine — Supervisor configuration
; Keeps the Python copy trading worker running and restarts it on failure.
;
; Install:
;   sudo cp python/copy_trading/supervisor/copyengine.conf /etc/supervisor/conf.d/
;   sudo supervisorctl reread && sudo supervisorctl update
;   sudo supervisorctl start copyengine
;
; Logs:
;   tail -f /var/log/tradesync/copyengine.log
; ─────────────────────────────────────────────────────────────────────────────

[program:copyengine]
command=python -m python.copy_trading.main
directory=/opt/tradesync
autostart=true
autorestart=true
startsecs=5
startretries=10
stopwaitsecs=30
stdout_logfile=/var/log/tradesync/copyengine.log
stderr_logfile=/var/log/tradesync/copyengine_err.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
user=tradesync
environment=PYTHONPATH="/opt/tradesync",PYTHONUNBUFFERED="1"
