# python/copy_trading/supervisor/mt5bridge.conf
; ─────────────────────────────────────────────────────────────────────────────
; MT5 Wine Bridge — Supervisor configuration
; Manages three processes in dependency order:
;   1. xvfb        — virtual display (MT5 needs a screen even on headless VPS)
;   2. mt5terminal — MetaTrader 5 terminal running inside Wine
;   3. mt5bridge   — mt5-remote RPC server (Linux Python connects here)
;
; Install:
;   sudo cp python/copy_trading/supervisor/mt5bridge.conf /etc/supervisor/conf.d/
;   sudo supervisorctl reread && sudo supervisorctl update && sudo supervisorctl start all
;
; Per-account scaling (multiple MT5 terminals):
;   Duplicate [program:mt5terminal] and [program:mt5bridge] blocks with
;   different --mt5path and --port values (18812, 18813, 18814…).
;   Update MT5_BRIDGE_PORT in .env for each worker that should use a
;   specific account port.
; ─────────────────────────────────────────────────────────────────────────────

; ── 1. Virtual display ────────────────────────────────────────────────────────
[program:xvfb]
command=Xvfb :99 -screen 0 1024x768x16
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/tradesync/xvfb.log
stderr_logfile=/var/log/tradesync/xvfb_err.log
user=root

; ── 2. MetaTrader 5 terminal (Wine) ──────────────────────────────────────────
[program:mt5terminal]
command=wine "%(ENV_HOME)s/.wine/drive_c/Program Files/MetaTrader 5/terminal64.exe"
environment=DISPLAY=":99",WINEPREFIX="%(ENV_HOME)s/.wine"
autostart=true
autorestart=true
priority=20
startsecs=10
stdout_logfile=/var/log/tradesync/mt5terminal.log
stderr_logfile=/var/log/tradesync/mt5terminal_err.log
user=root

; ── 3. mt5-remote bridge server (Wine Python) ────────────────────────────────
[program:mt5bridge]
command=wine python -m mt5_remote --host 0.0.0.0 --port 18812
environment=DISPLAY=":99",WINEPREFIX="%(ENV_HOME)s/.wine"
autostart=true
autorestart=true
priority=30
startsecs=15
stdout_logfile=/var/log/tradesync/mt5bridge.log
stderr_logfile=/var/log/tradesync/mt5bridge_err.log
user=root
